# USERHOME/kde_color/.config/systemd/user/kde_color.path
[Path]
PathChanged=~/.config/kdeglobals

[Install]
WantedBy=multi-user.target
