# USERHOME/kde_color/.config/systemd/user/kde_color.path
[Unit]
Requires=kde_color.service

[Path]
PathChanged=%h/.config/kdeglobals

[Install]
WantedBy=default.target
